# deps/perfetto/perfetto/debian/traced-probes.service
[Unit]
Description=Perfetto probes daemon

[Service]
ExecStart=/usr/sbin/traced_probes
User=root
Group=root
